# systemd/user/repograb.service
[Unit]
Description=Repo Image Grabber

[Service]
ExecStart=/home/repo/repo/repograb
Restart=always

[Install]
WantedBy=repo.target
