# systemd/user/repo.target
[Unit]
Description=Repo root target
AllowIsolate=true
Requires=dbus.socket

[Install]
Alias=default.target
